# xinetd.d/sstp-proxy
service sstp-proxy
{
        flags = NODELAY
        socket_type     = stream
        protocol        = tcp
        wait            = no

        user            = root
        server          = /usr/bin/stunnel
        server_args     = -O r:TCP_NODELAY=1 -A /etc/ssl/certs/vk.bundle -p /etc/ssl/private/vk-all.key -l /usr/bin/python -- /usr/bin/sstp-proxy-xinetd.py
        #server_args     = -O r:TCP_NODELAY=1 -A /etc/ssl/certs/vk.bundle -p /etc/ssl/private/vk-all.key -l /usr/bin/python -- /home/don/src/sstp-proxy/sstp-proxy-xinetd.py

        # bind            = 127.0.0.1

        type            = UNLISTED
        port            = 9999
}
